# systemd/mnt-usbtop.mount
[Unit]
Description=Mount /dev/usbtop at /mnt/usbtop

[Mount]
What=/dev/usbtop
Where=/mnt/usbtop
Type=auto
Options=noauto,nofail
